# deploy/dicty-app.service
[Unit]
Description=Dicty Resolwe Server Django Application
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service
Requires=docker.service

[Service]
Type=simple
User=dictyapp
Group=dictyapp
WorkingDirectory=/home/dictyapp/dicty_resolwe_server/backend/resolwe_server
Environment=PATH=/home/dictyapp/dicty_resolwe_server/backend/venv/bin
Environment=DJANGO_SETTINGS_MODULE=settings
ExecStartPre=/bin/sleep 10
ExecStart=/home/dictyapp/dicty_resolwe_server/backend/venv/bin/python manage.py runserver 0.0.0.0:8000
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dicty-app

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/dictyapp/dicty_resolwe_server
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
